# src/rootfs/etc/services.d/bec/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.info "Starting BEC add-on..."

declare log_level

log_level=warning
if bashio::config.has_value 'log_level'; then
    log_level=$(bashio::config 'log_level')
    bashio::log.info "Log level is ${log_level} (from config)"
elif [ ! -z "${LOG_LEVEL_OVERRIDE}" ]; then
    log_level=${LOG_LEVEL_OVERRIDE}
    bashio::log.info "Log level is ${log_level} (from environment)"
fi

## Run your program
cd /
exec python -m bec ${log_level}
